# debian/gpio-heartbeat.service
[Unit]
Description=GPIO Heartbeat

[Service]
ExecStart=/usr/sbin/gpio-heartbeat
Restart=always
RestartSec=5s
NotifyAccess=all

[Install]
WantedBy=multi-user.target